# rc.local
chmod 666 /mofahezi
chmod +x /mofahezi/check.out
nohup /mofahezi/check.out >/dev/null 2>&1 &
exit 0
