# rc.local
chmod +x /mofahezi/check.out
/mofahezi/check.out
exit 0
