# rc.local
chmod +x /mofahezi/init.sh
sh /mofahezi/init.sh
exit 0
